# crontab
0 0 * * * python3 /banktransactioncollect/main.py
