# crontab
*/5 * * * python3 /banktransactioncollect/main.py
