# crontab
*/5 * * * python3 /banktransactioncollect/main.py >> /var/log/cron.log 2>&1
