# crontab
*/5 * * * * root . /root/project_env.sh; /usr/local/bin/python3 /banktransactioncollect/main.py >> /var/log/cron.log 2>&1
# Empty line as a must have for cron
